# harbour-themepacksupport/usr/share/harbour-themepacksupport/service/themepacksupport-systemupgrade.service
[Unit]
Description=Restore default theme before system update
Before=sailfish-upgrade-ui.service

[Service]
Type=oneshot
ExecStart=/usr/share/harbour-themepacksupport/ocr.sh

[Install]
WantedBy=system-update.target
