# harbour-themepacksupport-0.8.0-6.noarch/usr/share/harbour-themepacksupport/service/harbour-themepacksupport.service
[Unit]
Description=Update icon theme via Theme pack support.

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/share/harbour-themepacksupport/service/run.sh

[Install]
WantedBy=multi-user.target
